# linux/systemd/phoneshell.service
[Unit]
Description=PhoneShell Headless
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/phoneshell
EnvironmentFile=-/etc/phoneshell/phoneshell.env
ExecStart=/opt/phoneshell/PhoneShell.Headless --config /etc/phoneshell/config.json
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
